# files/etc/systemd/system/systemd-nspawn@raspbian.service.d/override.conf
[Service]
DeviceAllow=char-video4linux rwm
DeviceAllow=char-vchiq rwm
